# files/geoipupdate.timer
[Unit]
Description=Daily GeoIP Update

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
