# files/reconta-notifications.service
[Unit]
Description=ReConta - Varredura de notificações de contas fixas
After=network.target reconta.service

[Service]
Type=oneshot
EnvironmentFile=-/var/www/reconta/api/.env
ExecStart=/bin/sh -c '/usr/bin/curl -fsS -X POST -H "X-Internal-Token: $INTERNAL_API_TOKEN" http://localhost:3020/api/internal/notifications/scan'
ExecStart=/bin/sh -c '/usr/bin/curl -fsS -X POST -H "X-Internal-Token: $INTERNAL_API_TOKEN" http://localhost:3020/api/internal/billing/scan'
