# files/reconta-notifications.timer
[Unit]
Description=Dispara a varredura de notificações de contas fixas a cada hora

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
