# files/reconta-cron.timer
[Unit]
Description=ReConta — Timer diário para notificações de contas a vencer

[Timer]
# Executa todo dia às 08:00
OnCalendar=*-*-* 08:00:00
# Garante reexecução caso o sistema esteja desligado no horário agendado
Persistent=true
# Aleatoriza até 5 minutos para evitar avalanche de execuções simultâneas
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
